# rootfs/etc/service/h2o/run
#!/bin/sh

exec /usr/local/bin/h2o -c /etc/h2o.conf
